# scripts/clock.service
# The startup script for the clock.
# This should go into /etc/systemd/system/

# Don't forget to modify WorkingDirectory and ExecStart
# to point to the correct directory and file.

[Unit]
Description = Book Quote Clock
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/Literary-Quote-Clock/
ExecStart=sudo python3 /path/to/Literary-Quote-Clock/clock.py
KillSignal=SIGINT
Restart=on-failure

[Install]
WantedBy=multi-user.target